# deploy/fanout-live.service
[Unit]
Description=Fanout Live RTMP Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/fanout-live
Environment=TWITCH_STREAM_KEY=replace-me
Environment=YOUTUBE_STREAM_KEY=replace-me
ExecStart=/usr/bin/python3 -m fanout_live --web --config /opt/fanout-live/config.toml --web-host 0.0.0.0 --web-port 8080
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
